# supervisor-vpn-poller.conf
[program:vpn-poller-realtime]
process_name=%(program_name)s
command=php /var/www/aiovpn/artisan vpn:poll-server --interval=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/aiovpn/storage/logs/vpn-poller.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
